# contrib/init/greencoind.conf
description "GreenCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GREENCOIND_BIN="/usr/bin/greencoind"
env GREENCOIND_USER="greencoin"
env GREENCOIND_GROUP="greencoin"
env GREENCOIND_PIDDIR="/var/run/greencoind"
# upstart can't handle variables constructed with other variables
env GREENCOIND_PIDFILE="/var/run/greencoind/greencoind.pid"
env GREENCOIND_CONFIGFILE="/etc/greencoin/greencoin.conf"
env GREENCOIND_DATADIR="/var/lib/greencoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # greencoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GREENCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run greencoind."
        echo "The setting must appear in $GREENCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"GreenCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GREENCOIND_PIDDIR"
    chmod 0755 "$GREENCOIND_PIDDIR"
    chown $GREENCOIND_USER:$GREENCOIND_GROUP "$GREENCOIND_PIDDIR"
    chown $GREENCOIND_USER:$GREENCOIND_GROUP "$GREENCOIND_CONFIGFILE"
    chmod 0660 "$GREENCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GREENCOIND_PIDFILE" \
    --chuid $GREENCOIND_USER:$GREENCOIND_GROUP \
    --exec "$GREENCOIND_BIN" \
    -- \
    -pid="$GREENCOIND_PIDFILE" \
    -conf="$GREENCOIND_CONFIGFILE" \
    -datadir="$GREENCOIND_DATADIR" \
    -disablewallet \
    -daemon

